# etc/systemd/system/peppymeterbasic.service.d/10-groups.conf
[Service]
SupplementaryGroups=render video
Environment=MESA_SHADER_CACHE_DIR=/tmp/mesa_shader_cache
